# app-admin/guix/files/guixd-0.8
#!/sbin/runscript

DAEMON=/usr/bin/guix-daemon
DAEMON_OPTS="--build-users-group=guix-builders --cores=0"
NAME=guixd
DESC=guix-daemon

depend() {
    	need localmount
}


start() {
    ebegin "Starting $DESC: "
    start-stop-daemon --start --background --make-pidfile \
	--quiet --pidfile /var/run/$NAME.pid \
	--exec $DAEMON -- $DAEMON_OPTS
    eend $? "Starting $DESC failed"
}

stop() {
    ebegin "Stopping $DESC: "
    start-stop-daemon --stop --quiet --pidfile /var/run/$NAME.pid \
	--exec $DAEMON
    eend $? "Stopping $DESC failed"
}

